# scripts/systemd/highfive-deploy.timer
# Fires highfive-deploy.service ~every 2 minutes. OnUnitInactiveSec measures from
# the previous run's COMPLETION, so a long build never queues overlapping ticks
# (flock in deploy.sh is the second guard). The cheap "anything new on main?"
# check exits in milliseconds when main is unchanged.
[Unit]
Description=Run highfive-deploy every 2 minutes

[Timer]
OnBootSec=2min
OnUnitInactiveSec=2min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
